# mid-service/asset/supervisord.conf
[supervisord]
nodaemon = true
user = root

[program:mid]
command = /opt/agent/bin/mid.sh start
stdout_logfile=/opt/agent/logs/agent0.log.0
stdout_logfile_maxbytes=0
user = root